# docker/crontab
# m h  dom mon dow   command
15 14 * * * /usr/bin/python3 /webapp/etl_cron.py >> /webapp/cron.log 2>&1
# remember to end this file with an empty new line